# docker/queue-worker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work rabbitmq --queue=default --sleep=3 --tries=3 --max-time=3600
stdout_logfile=/var/www/html/storage/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
user=www-data
autostart=true
autorestart=true
numprocs=8
redirect_stderr=true
stopwaitsecs=3600

